# backend/services/weltenwind-backend.service
[Unit]
Description=Weltenwind Backend API Service
Documentation=https://github.com/weltenwind/backend
After=network.target postgresql.service
PartOf=weltenwind.target

[Service]
Type=simple
User=weltenwind
Group=weltenwind
WorkingDirectory=/srv/weltenwind/backend
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="LOG_LEVEL=info"
EnvironmentFile=/srv/weltenwind/backend/.env
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10

# Logging Configuration
# systemd Standard-Logs (stdout/stderr)
StandardOutput=append:/var/log/weltenwind/backend.log
StandardError=append:/var/log/weltenwind/backend.error.log

# Winston structured logs gehen nach /var/log/weltenwind/
# (app.log, auth.log, security.log, api.log, error.log)

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/weltenwind /srv/weltenwind/backend

[Install]
WantedBy=weltenwind.target multi-user.target 