# backend/services/weltenwind-studio.service
[Unit]
Description=Weltenwind Prisma Studio Service
Documentation=https://github.com/weltenwind/backend
After=network.target postgresql.service
PartOf=weltenwind.target

[Service]
Type=simple
User=weltenwind
Group=weltenwind
WorkingDirectory=/srv/weltenwind/backend
Environment="NODE_ENV=production"
Environment="HOME=/var/lib/weltenwind"
EnvironmentFile=/srv/weltenwind/backend/.env
# Stelle sicher dass npm/node im PATH ist
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/bin/bash -c 'cd /srv/weltenwind/backend && test -f node_modules/.bin/prisma || (echo "Prisma not installed!" && exit 1)'
ExecStart=/bin/bash -c 'cd /srv/weltenwind/backend && /usr/bin/npx prisma studio --port 5555 --browser none --hostname 0.0.0.0'
# Health Check - Port 5555 muss erreichbar sein
ExecStartPost=/bin/bash -c 'for i in {1..20}; do nc -z localhost 5555 >/dev/null 2>&1 && break || sleep 3; done'
# Advanced Restart Policy
Restart=always
RestartSec=15
RestartSteps=5
RestartMaxDelaySec=180
StandardOutput=append:/var/log/weltenwind/studio.log
StandardError=append:/var/log/weltenwind/studio.error.log

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/weltenwind

[Install]
WantedBy=weltenwind.target multi-user.target 