# backend/services/weltenwind-backend.service
[Unit]
Description=Weltenwind Backend API Service
Documentation=https://github.com/weltenwind/backend
After=network.target postgresql.service
PartOf=weltenwind.target

[Service]
Type=simple
User=weltenwind
Group=weltenwind
WorkingDirectory=/srv/weltenwind/backend
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="LOG_LEVEL=info"
EnvironmentFile=/srv/weltenwind/backend/.env
ExecStart=/usr/bin/npm run dev
# Health Check - API muss erreichbar sein
ExecStartPost=/bin/bash -c 'for i in {1..30}; do curl -f http://localhost:3000/api/health >/dev/null 2>&1 && break || sleep 2; done'
# Advanced Restart Policy
Restart=always
RestartSec=10
# Exponential Backoff für kritische Failures
RestartSteps=10
RestartMaxDelaySec=300
# Watchdog für Health Monitoring (alle 30s)
WatchdogSec=30

# Logging Configuration
# systemd Standard-Logs (stdout/stderr)
StandardOutput=append:/var/log/weltenwind/backend.log
StandardError=append:/var/log/weltenwind/backend.error.log

# Winston structured logs gehen nach /var/log/weltenwind/
# (app.log, auth.log, security.log, api.log, error.log)

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/weltenwind /srv/weltenwind/backend

[Install]
WantedBy=weltenwind.target multi-user.target 